# software/ansible/roles/pcf8523/files/85-hwclock-pcf8523.rules
## Generated by Ansible (pcf8523)
#  (run 'udevadm trigger' to enforce changes)

# Set the system time as soon as the hardware clock shows up
KERNEL=="rtc0", RUN+="/etc/udev/hwclock-pcf8523-set $root/$name"
