# software/ansible/roles/ntp/files/ntp-wait.service
## Generated by Ansible (ntp)

# Notify systemd-time-wait-sync.service when NTP has synchronized system time
# (see ntp-wait manpage)
[Unit]
Description=Wait for NTP time synchronization
Requires=ntp.service
After=ntp.service

[Service]
Type=oneshot
ExecStartPre=mkdir -p /run/systemd/timesync
ExecStart=/usr/sbin/ntpwait -n 600 -s 6
ExecStartPost=touch /run/systemd/timesync/synchronized
TimeoutStartSec=infinity
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
