# software/ansible/roles/reboot/templates/scheduled-reboot.service.j2
## Generated by Ansible (reboot)

[Unit]
Description=Scheduled Reboot

[Service]
Type=simple
StandardOutput=journal+console
ExecStart=systemctl reboot{% if REBOOT_SCHEDULED_FORCED|bool %} --force{% if REBOOT_SCHEDULED_UNGRACEFUL|bool %} --force{% endif %}{% endif %}

[Install]
# -> scheduled-reboot.timer
