# software/ansible/roles/reboot/templates/scheduled-reboot.timer.j2
## Generated by Ansible (reboot)

[Unit]
Description=Scheduled Reboot

[Timer]
OnCalendar={{ REBOOT_SCHEDULED_ONCALENDAR }}
Unit=scheduled-reboot.service

[Install]
WantedBy=timers.target
