# software/ansible/roles/ethernet/files/ethernet-power-down.service
## Generated by Ansible (ethernet)

[Unit]
Description=Ethernet Power-Down
After=network.target

[Service]
Type=simple
StandardOutput=journal+console
ExecStart=ip link set dev eth0 down

[Install]
WantedBy=multi-user.target
