# software/ansible/roles/sleepypi/files/sleepy-pi-watchdog.service
## Generated by Ansible (sleepypi)

[Unit]
Description=Sleepy Pi Watchdog

[Service]
Type=simple
StandardOutput=journal+console
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/default/sleepy-pi-watchdog
ExecStart=/usr/local/bin/sleepy-pi-watchdog

[Install]
WantedBy=multi-user.target
