# software/ansible/roles/sleepypi/files/sleepy-pi-agent.service
## Generated by Ansible (sleepypi)

[Unit]
Description=Sleepy Pi Agent

[Service]
Type=simple
StandardOutput=journal+console
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/local/bin/sleepy-pi-agent

[Install]
WantedBy=multi-user.target
