# software/ansible/roles/sleepypi/files/sleepy-pi-sunwatch.service
## Generated by Ansible (sleepypi)

[Unit]
Description=Sleepy Pi Sunwatch

[Service]
Type=simple
StandardOutput=journal+console
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/etc/default/sleepy-pi-sunwatch
ExecStart=/usr/local/bin/sleepy-pi-sunwatch

[Install]
WantedBy=multi-user.target
