# software/ansible/roles/ogn/files/ogn-rf-fifo.service
## Generated by Ansible (ogn)

[Unit]
Description=Open Glider Network (OGN) Receiver (RF) FIFO file
ConditionPathExists=!/run/ogn-rf.fifo
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=mkfifo /run/ogn-rf.fifo
ExecStart=chgrp _ogn /run/ogn-rf.fifo
ExecStart=chmod 660 /run/ogn-rf.fifo
ExecStop=rm -f /run/ogn-rf.fifo
