# software/ansible/roles/freedns/files/freedns@.service
## Generated by Ansible (freedns)

## FreeDNS client service
[Unit]
Description=FreeDNS client (%i)

[Service]
Type=simple
StandardOutput=journal+console
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/local/bin/freedns-v2 --config /etc/freedns/%i.conf --foreground

[Install]
WantedBy=network.target
